# systemd-user/qrexec-connect.service
[Unit]
Description=systemd-native qrexec-client-vm service

[Service]
Type=notify
StandardError=journal
ExecStart=/opt/bin/qrexec-connect
